# playbooks/files/nat-routing.service
[Unit]
Description=NAT Routing Rules
After=network-pre.target
Before=network.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/sbin/iptables-restore /etc/iptables/nat-routing.rules
ExecStop=/sbin/iptables -t nat -F
ExecStop=/sbin/iptables -t filter -F
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
