# playbooks/files/network_prober.service
[Unit]
Description=Network prober service

[Service]
ExecStart=/home/pi/network_probe_runner.sh
Restart=always

